# rootfs/systemd/hello.service
[Unit]
Description=Hello service for CrisOS
WantedBy=default.target

[Service]
ExecStart=echo Hello from systemd service
