# rootfs/openrc/hello.service
[Unit]
Description=Hello service for NucleOS
WantedBy=default.target

[Service]
ExecStart=echo Hello from openrc service
